# images/amavis/rootfs/etc/services.d/20-cron/run
#!/bin/sh

exec crond -f -S -l $CRON_LOG_LEVEL
