# service/clamav/rootfs/etc/services.d/10-clamd-log-forwarder/run
#!/bin/sh

LOGFILE=/var/log/clamav/clamd.log

[ -e $LOGFILE ] \
  || mkfifo $LOGFILE \
  && chown clamav:clamav $LOGFILE

exec tail -f -n 0 $LOGFILE
