# images/amavis/rootfs/etc/periodic/daily/razor-admin-discover
#!/bin/sh

razor-admin $RAZOR_OPTS -home=/var/amavis/.razor -discover
chown -R amavis:amavis /var/amavis/.razor
