# images/amavis/rootfs/etc/periodic/daily/spamassassin-update
#!/bin/sh

/usr/bin/sa-update $SA_UPDATES_LOG_LEVEL_FLAG
